# cron/crontab
* * */1 * * /usr/local/bin/python /sample_crawler/trend.py > /proc/1/fd/1 2>/proc/1/fd/2
